# dots/home/aiden/.config/systemd/user/timers.target.wants/swww.timer
[Unit]
Description=Refreshes swww wallpaper

OnFailure=notify-err@%n

[Timer]
OnUnitActiveSec=300

[Install]
WantedBy=timers.target
